# dijkstra-db.service
[Unit]
Description=Dijkstra DB Service

[Service]
User=app
ExecStart=/usr/bin/docker start mysql-product

[Install]
WantedBy=multi-user.target
